# deployment/marketscalper-backup.timer
# MarketScalper daily backup schedule (Phase D).
# Enable: sudo systemctl enable --now marketscalper-backup.timer
# Check:  systemctl list-timers marketscalper-backup.timer
[Unit]
Description=Daily MarketScalper database backup

[Timer]
# 03:15 UTC every day (explicit UTC so the VPS location — e.g. Canada — never
# shifts it); catch up if the machine was off.
OnCalendar=*-*-* 03:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
