# deployment/marketscalper.service
# MarketScalper systemd unit (roadmap P0.27; Decision D4).
# Single Linux server, single process, provider-agnostic.
# Install: copied to /etc/systemd/system/ by deployment/deploy.sh.

[Unit]
Description=MarketScalper — deterministic market analysis & decision-support terminal
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=marketscalper
Group=marketscalper
WorkingDirectory=/opt/marketscalper
# Secrets (DB DSN, API token) live only here — never in git (D3/D4):
EnvironmentFile=/etc/marketscalper/env
ExecStart=/opt/marketscalper/venv/bin/marketscalper
Restart=always
RestartSec=5

# --- sandbox hardening (roadmap P4.13). Decision-support only: the
# process needs outbound network (feed) + the DB, its read-only code
# tree, and its own log dir — nothing else.
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
# ProtectSystem=strict makes the FS read-only; only the log dir is writable.
ReadWritePaths=/var/log/marketscalper
# Resource guards for the single-process VPS.
LimitNOFILE=8192
TasksMax=256

[Install]
WantedBy=multi-user.target
